# deploy/micro_credit.service
[Unit]
Description=Micro Credit Backend (Docker)
After=network.target

[Service]
User=root
WorkingDirectory=/opt/micro_credit
EnvironmentFile=/opt/micro_credit/.env
Restart=always
RestartSec=10
# If you're using docker-compose, start it here. Adjust paths as needed.
ExecStart=/usr/bin/docker compose up --force-recreate
ExecStop=/usr/bin/docker compose down

[Install]
WantedBy=multi-user.target
